# home/.config/autostart/copyq.desktop
[Desktop Entry]
Type=Application
Name=CopyQ
Exec=copyq
X-GNOME-Autostart-enabled=true
NoDisplay=false
StartupNotify=false
Categories=Utility;
